# files/local-dehydrated.service
[Unit]
Description=dehydrated
[Service]
Type=oneshot
ExecStart=/usr/bin/dehydrated -c
PrivateTmp=true
